# deployment/supervisor/reverb.conf
[program:reverb]
process_name=%(program_name)s
command=php /var/www/gestao_frota/artisan reverb:start
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/gestao_frota/storage/logs/reverb.log
stopwaitsecs=3600
